# git-server/local-git-daemon.conf
start on startup
stop on shutdown
exec /usr/bin/git daemon \
    --user=go --group=go \
    --export-all \
    --reuseaddr \
    --base-path=/var/go/pipelines \
    /var/go/pipelines
respawn
